# gentoo/etc/init.d/hv_kvp_daemon
#!/sbin/openrc-run

# Copyright 2019 Christopher Dangerfield <DL200010@gmail.com>
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.


# Do this to get the hv_kvp_daemon
# cd /usr/src/linux/
# make
# cd /usr/src/linux/tools/hv/
# make install

command="/usr/sbin/hv_kvp_daemon"
pidfile="/run/hv_kvp_daemon.pid"
name="hv_kvp_daemon"

depend() {
    need net
}

start_pre() {
    checkpath --directory --mode 0755 /run
}

start() {
    ebegin "Starting ${name}"
    start-stop-daemon --start --pidfile ${pidfile} --exec ${command}
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    start-stop-daemon --stop --pidfile ${pidfile}
    eend $?
}
